# ops/systemd/colibri-watchdog.service
[Unit]
Description=Check Colibri for an idle inference wedge
After=colibri-server.service

[Service]
Type=oneshot
EnvironmentFile=%h/.config/colibri/service.env
ExecStart=%h/colibri/ops/colibri-watchdog.sh
