# ops/systemd/colibri-server.service.d/10-preflight.conf
# Drop-in: refuse to start before the host is ready, instead of failing deep
# inside a multi-minute model load.  Enabled by the COLI_PREFLIGHT_* variables
# in service.env (see ops/service.env.multigpu.example); with none of them set
# the script only verifies that the model directory looks like a snapshot.
#
# RestartSec here is deliberately long: after a reboot the tmpfs copy takes
# ~15 minutes, and a 10-second retry loop would just compete with it for I/O.
[Service]
ExecStartPre=%h/colibri/ops/colibri-preflight.sh
RestartSec=60
