# ops/systemd/colibri-server.service.d/50-engine-openmp.conf
# Keep the native HTTP mux unbound until it forks the inference engine.  The
# engine launcher then restores the tuned binding policy only in the child.
# Setting OMP_PROC_BIND=spread on the mux itself lets libgomp pin it before
# fork(), causing every engine worker to inherit a single-core affinity mask.
[Service]
Environment=OMP_PROC_BIND=false
Environment=COLI_ENGINE_OMP_PROC_BIND=spread
